# Lecture-38/crontab
* * * * * /usr/local/bin/get_metrics.py
